# overlays/xpra/etc/profile.d/20-amazonspiceox-xpra.sh
export AMAZONSPICEOX_XPRA_PROFILE=1
export ASOX_GUI_BACKEND="${ASOX_GUI_BACKEND:-xpra}"
export XPRA_PORT="${XPRA_PORT:-14500}"
export XPRA_DISPLAY_NUM="${XPRA_DISPLAY_NUM:-100}"
