# overlays/gui/etc/profile.d/20-amazonspiceox-gui.sh
export AMAZONSPICEOX_GUI_PROFILE=1
export BROWSER=/usr/local/bin/asox-browser
export ASOX_GUI_BACKEND="${ASOX_GUI_BACKEND:-auto}"
export DISPLAY="${DISPLAY:-10.0.2.2:0}"
export LIBGL_ALWAYS_SOFTWARE=1
export NO_AT_BRIDGE=1
