# overlays/ssm-powerconnect/etc/profile.d/20-amazonspiceox-ssm-powerconnect.sh
export AMAZONSPICEOX_SSM_POWERCONNECT_PROFILE=1
export SSM_POWERCONNECT_HOME="${SSM_POWERCONNECT_HOME:-/opt/ssm-powerconnect}"
