# overlays/terraform/etc/profile.d/20-amazonspiceox-terraform.sh
export AMAZONSPICEOX_TERRAFORM_PROFILE=1
