# overlays/docker/etc/profile.d/20-amazonspiceox-docker.sh
export AMAZONSPICEOX_DOCKER_PROFILE=1
export DOCKER_HOST="${DOCKER_HOST:-unix:///var/run/docker.sock}"
