# overlays/kubectl/etc/profile.d/20-amazonspiceox-kubectl.sh
export AMAZONSPICEOX_KUBECTL_PROFILE=1
export KUBECONFIG=/root/.kube/config
