# overlays/debug/etc/profile.d/20-amazonspiceox-debug.sh
export AMAZONSPICEOX_DEBUG_PROFILE=1
export PS1='arrakis[debug]:\w# '
